# linux/systemd/myhome-update.timer
[Unit]
Description=Daily update of myhome package

[Timer]
OnCalendar=daily
AccuracySec=1h

[Install]
WantedBy=timers.target
